# sfwebsite.service
[Unit]
Description=SF Learning Platform (FastAPI + vscoqtop)
After=network.target

[Service]
Type=exec
User=www-data
Group=www-data
WorkingDirectory=/var/www/sfwebsite
EnvironmentFile=/var/www/sfwebsite/.env
ExecStart=/var/www/sfwebsite/venv/bin/uvicorn server.main:app --host 127.0.0.1 --port 8100
Restart=always
RestartSec=5
# Limit memory to prevent runaway vscoqtop processes
MemoryMax=2G

[Install]
WantedBy=multi-user.target
